# kovh/data/systemd/flanneld.service.d/10-client.conf
[Service]
Environment="FLANNELD_ETCD_ENDPOINTS=http://192.168.0.1:2379"
