# kovh/data/systemd/flanneld.service.d/10-client.conf
[Service]
Environment="RKT_RUN_ARGS=--uuid-file-save=/var/lib/coreos/flannel-wrapper.uuid --volume tls-certs,kind=host,source=/etc/kubernetes/tls,readOnly=true --mount volume=tls-certs,target=/tls"
Environment="FLANNELD_ETCD_ENDPOINTS=https://localhost:2379"
Environment="FLANNELD_ETCD_CAFILE=/tls/ca.pem"
Environment="FLANNELD_ETCD_CERTFILE=/tls/etcd_c.crt"
Environment="FLANNELD_ETCD_KEYFILE=/tls/etcd_c.key"
