# kovh/data/systemd/flanneld.service.d/10-client.conf
[Service]
Environment="FLANNELD_ETCD_ENDPOINTS=http://localhost:2379"
