# kovh/data/systemd/flanneld.service.d/10-daemon.conf
[Service]
Environment="RKT_RUN_ARGS=--uuid-file-save=/var/lib/coreos/flannel-wrapper.uuid --volume tls-certs,kind=host,source=/etc/kubernetes/tls,readOnly=true --mount volume=tls-certs,target=/tls"
Environment="FLANNELD_IFACE=eth1"
Environment="FLANNELD_ETCD_ENDPOINTS=__ETCD_URL__"
Environment="FLANNELD_ETCD_CAFILE=/tls/ca.pem"
Environment="FLANNELD_ETCD_CERTFILE=/tls/etcdclient.crt"
Environment="FLANNELD_ETCD_KEYFILE=/tls/client.key"
