# kovh/data/systemd/flanneld.service.d/10-daemon.conf
[Service]
Environment="RKT_RUN_ARGS=--uuid-file-save=/var/lib/coreos/flannel-wrapper.uuid --volume tls-certs,kind=host,source=/etc/kubernetes/tls,readOnly=true --mount volume=tls-certs,target=/etc/kubernetes/tls"
Environment="FLANNELD_IFACE=eth1"
Environment="FLANNELD_ETCD_ENDPOINTS=__ETCD_URL__"
Environment="FLANNELD_ETCD_CAFILE=/etc/kubernetes/tls/ca.pem"
Environment="FLANNELD_ETCD_CERTFILE=/etc/kubernetes/tls/client/etcd.crt"
Environment="FLANNELD_ETCD_KEYFILE=/etc/kubernetes/tls/host.key"
