# kovh/data/systemd/flanneld.service.d/10-network-config.conf
[Service]
ExecStartPre=/usr/bin/etcdctl \
  --endpoints https://localhost:2379 \
  --ca-file /etc/kubernetes/tls/ca.pem \
  --cert-file /etc/kubernetes/tls/client/etcd.crt \
  --key-file /etc/kubernetes/tls/host.key \
  set /coreos.com/network/config '{ "Network": "172.17.0.0/16", "Backend": { "Type": "vxlan" } }'
