# kovh/data/systemd/coreos-metadata-sshkeys@.service.d/10-provider.conf
[Service]
ExecStart=
ExecStart=/usr/bin/coreos-metadata \
  --provider=openstack-metadata \
  --ssh-keys=%i
