# kovh/data/systemd/etcd-member.service.d/10-cluster.conf
[Unit]
Requires=coreos-metadata.service
After=coreos-metadata.service

[Service]
EnvironmentFile=/run/metadata/coreos
ExecStart=
ExecStart=/usr/lib/coreos/etcd-wrapper \
  -listen-client-urls http://${COREOS_OPENSTACK_IPV4_LOCAL}:2379,http://localhost:2379 \
  -listen-peer-urls http://${COREOS_OPENSTACK_IPV4_LOCAL}:2380,http://localhost:2380 \
  -advertise-client-urls http://${COREOS_OPENSTACK_IPV4_LOCAL}:2379 \
  -initial-advertise-peer-urls http://${COREOS_OPENSTACK_IPV4_LOCAL}:2380 \
  -initial-cluster ${ETCD_NAME}=http://${COREOS_OPENSTACK_IPV4_LOCAL}:2380
